# Daily_Regression_Test/tele_op_services/com.prismax.teleop.regression.catchup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.prismax.teleop.regression.catchup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/wanxin/PycharmProjects/WORK/Prismax/QA_PrismaX/Daily_Regression_Test/tele_op_services/run_tests_if_needed.sh</string>
    </array>

    <!-- 登录时运行一次 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 每 30 分钟检查一次（脚本内部有时间窗口保护，只在 09:00–11:00 真正触发） -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/Users/wanxin/PycharmProjects/WORK/Prismax/QA_PrismaX/Daily_Regression_Test/tele_op_services/log/regression_catchup_launchagent.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/wanxin/PycharmProjects/WORK/Prismax/QA_PrismaX/Daily_Regression_Test/tele_op_services/log/regression_catchup_launchagent.log</string>
</dict>
</plist>
