# Daily_Pull/com.prismax.daily_pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.prismax.daily_pull</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>bash</string>
        <string>/Users/wanxin/PycharmProjects/Prismax/QA_PrismaX/Daily_Pull/run_daily_pull_once.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/wanxin/PycharmProjects/Prismax/QA_PrismaX/Daily_Pull</string>
    <!-- 每天 9:00 若已登录则执行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <!-- 兜底：每 60 分钟检查一次，由脚本内部的 .last_run_date 逻辑保证同一天只会真正执行一次 -->
    <key>StartInterval</key>
    <integer>3600</integer>
    <!-- 登录/首次加载本 agent 时执行一次（登录时），用于补跑（若当天未跑过则跑） -->
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/wanxin/PycharmProjects/Prismax/QA_PrismaX/Daily_Pull/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/wanxin/PycharmProjects/Prismax/QA_PrismaX/Daily_Pull/launchd_stderr.log</string>
</dict>
</plist>
